# crontab
# m h DoM Mon DoW command
* * * * * /bin/bash /memorious/runcrawlers.sh 1> /var/log/cron.log 2> /var/log/cron.log
