# crontab
# m h DoM Mon DoW command
0 0 * * * /bin/bash /memorious/runcrawlers.sh >> /var/log/cron.log 2>&1
