# crontab
# m h DoM Mon DoW command
* * * * * /bin/bash /memorious/runcrawlers.sh >> /var/log/cron.log 2>&1
